# vm-images/android/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1920x1080x24 -ac +extension GLX +render -noreset
autostart=true
autorestart=true
priority=1
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb.error.log

[program:fluxbox]
command=/usr/bin/fluxbox
environment=DISPLAY=:99
autostart=true
autorestart=true
priority=2
stdout_logfile=/var/log/supervisor/fluxbox.log
stderr_logfile=/var/log/supervisor/fluxbox.error.log

[program:vm-host]
command=node /app/dist/index.js
directory=/app
autostart=true
autorestart=true
priority=3
stdout_logfile=/var/log/supervisor/vm-host.log
stderr_logfile=/var/log/supervisor/vm-host.error.log
